# config/tangle.conf
#!upstart
description 'Genie Tangle'

env USER=passenger
env RACK_ENV=production
env DAEMON=/usr/local/rbenv/shims/bundle

start on (started network-interface
          or started network-manager
          or started networking)

stop on (started network-interface
         or started network-manager
         or started networking)

respawn

exec start-stop-daemon --start --make-pidfile      \
  --pidfile '<%= shared_path %>/pids/tangle.pid'   \
  --chuid $USER                                    \
  --chdir '<%= current_path %>'                    \
  --exec $DAEMON                                   \
  -- exec tangle server                            \
    --port=<%= app_port %>                         \
    --faye-port=<%= faye_port %>                   \
    --log-file='<%= shared_path %>/log/production.log'

respawn limit 10 90
